# kerberos-files/etc/profile.d/graphics.sh
[ -x /opt/pixar ] && export RMANTREE="/opt/pixar/RenderManProServer-21.8"
[ -x /opt/pixar ] && export MAYA_DISABLE_CIP=1
