# kerberos-files/etc/profile.d/graphics.csh
if ( -x /opt/pixar ) setenv RMANTREE "/opt/pixar/RenderManProServer-21.8A"
if ( -x /opt/pixar ) setenv MAYA_DISABLE_CIP 1
