# kerberos-files/etc/cron.d/checkgssd
30 * * * * root /usr/libexec/checkgssd >/tmp/checkgssd.log 2>&1

