# kerberos-files/etc/cron.d/killjob
12 * * * * root /usr/libexec/killjob.sh >/dev/null 2>&1
