# kerberos-files/etc/cron.d/killjob~
12 * * * * root /usr/libexec/killlong.sh >/dev/null 2>&1
