# kerberos-files/etc/systemd/system/sssd.service.d/override.conf
[Unit]
Requires=systemd-resolved.service
After=systemd-resolved.service
